# documents/workers/jobs/scripts/services/demo/ktz.ini
[program:ktz]
autorestart=True
autostart=True
redirect_stderr=True
command=/usr/local/bin/python3 /data/ktz/manage.py runserver 0.0.0.0:9092
user=root
directory=/data/ktz
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /data/ktz/logs/ktz.log